# ccdc_2024/AnsibleScripts/playbooks/firewall/roles/blackout/files/iptables-persistent
#!/sbin/openrc-run

depend() {
  need net
}

command="/bin/bash"
command_args="/etc/iptables/restore-iptables.sh"
pidfile="iptables-persistent.pid"
